# .config/autostart/killall gnome-documents.desktop
[Desktop Entry]
Name=killall gnome-documents
Comment=Run killall gnome-documents
Exec=killall gnome-documents
Terminal=false
Type=Application